# puppet/modules/faro_api/supervisor/faro_api.conf
[program:faro_api]
command = /opt/faro/faro_api/service-start.sh
autostart = true
autorestart = true
user = faro
stdout_logfile = /opt/faro/faro_api/logs/faro_api.log
redirect_stderr = true
